# deploy/price-tracker-bot.service
[Unit]
Description=Bot de Telegram rastreador de precios
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Bot_Telegram_Precio
Environment=PATH=/root/Bot_Telegram_Precio/venv/bin:/usr/bin:/bin
ExecStart=/root/Bot_Telegram_Precio/venv/bin/python main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
